# files/etc/uci-defaults/88-remove-unavailable-repo
sed -i "/sundaqiang/d" /etc/apk/repositories.d/distfeeds.list
sed -i "/fancontrol/d" /etc/apk/repositories.d/distfeeds.list
